# etc/cron.hourly/hpcc_certificates
#!/bin/bash

cd /usr/local/bin
./copy_lets_encrypt_certs.sh

